# systemd/picframe-lights.service
[Unit]
Description=PicFrame Indicator Lights Overlay
After=picframe-api.service
Wants=picframe-api.service

[Service]
Type=simple
Environment=DISPLAY=:0
Environment=GDK_BACKEND=x11
ExecStart=/usr/bin/python3 %h/picframe_4.0/scripts/indicator_lights.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
