# systemd/picframe-sync.timer
[Unit]
Description=PicFrame 4.0 Photo Sync Timer

[Timer]
# Run 15 minutes after boot
OnBootSec=15min
# Fire at clock-aligned 15-minute marks (:00, :15, :30, :45)
# Updated dynamically by /api/settings when sync interval changes
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
