# scripts/setup/systemd/picframe-watchdog.service
[Unit]
Description=PicFrame WiFi Watchdog
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/home/matt/picframe_4.0/scripts/setup/venv/bin/python3 /home/matt/picframe_4.0/scripts/setup/watchdog.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
