# scripts/setup/systemd/picframe-ap-setup.service
[Unit]
Description=PicFrame AP Setup Mode (hotspot + captive portal)
After=network.target
# Do not run automatically — started by picframe-watchdog when needs_setup=true

[Service]
Type=simple
ExecStartPre=/usr/sbin/hostapd -B /etc/hostapd/picframe-hostapd.conf
ExecStartPre=/usr/sbin/dnsmasq -C /tmp/picframe-dnsmasq.conf --pid-file=/tmp/picframe-dnsmasq.pid
ExecStartPre=/sbin/ip addr add 192.168.4.1/24 dev wlan0
ExecStart=/home/matt/picframe_4.0/scripts/setup/venv/bin/python3 /home/matt/picframe_4.0/scripts/setup/ap_portal/portal.py
ExecStopPost=/bin/sh -c 'kill $(cat /tmp/picframe-dnsmasq.pid) 2>/dev/null || true'
ExecStopPost=/usr/bin/pkill -f picframe-hostapd.conf || true
ExecStopPost=/sbin/ip addr del 192.168.4.1/24 dev wlan0 || true
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
# Not enabled by default — started on demand by watchdog
WantedBy=multi-user.target
