# scripts/setup/systemd/picframe-ble-setup.service
[Unit]
Description=PicFrame BLE Setup Service
After=bluetooth.target
Requires=bluetooth.target
# Do not run automatically — started by picframe-watchdog when needs_setup=true

[Service]
Type=simple
# __PROJECT_DIR__ is substituted with the actual project path by install_setup.sh
ExecStart=__PROJECT_DIR__/scripts/setup/venv/bin/python3 __PROJECT_DIR__/scripts/setup/ble_setup.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
# Not enabled by default — started on demand by watchdog
WantedBy=multi-user.target
